# gollum/gollum.socket
[Unit]
Description=Gollum Git-powered Wiki Socket

[Socket]
ListenStream=0.0.0.0:80

[Install]
WantedBy=sockets.target
